# .archive/legacy_docs/v4_structure/SHawn_Brain/shawn-bot.service
[Unit]
Description=SHawn-Lab Bot Auto-Service (Immortal)
After=network.target

[Service]
Type=simple
User=leseichi
WorkingDirectory=/home/leseichi/SHawn-BOT
# auto_restart.sh를 실행하여 git pull & 무한 재시작 관리
ExecStart=/bin/bash /home/leseichi/SHawn-BOT/auto_restart.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
